# home/dot_config/systemd/user/absurd-cleanup.timer
[Unit]
Description=Run Absurd cleanup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
