# home/dot_config/systemd/user/vicinae-cleanup.timer
[Unit]
Description=Run vicinae-cleanup hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
